# services/sync/crontab
*/10 * * * * root /scripts/export_notes.sh >> /var/log/cron.log 2>&1
