# services/delivery/crontab
*/5 * * * * root python3 /App/main.py >> /var/log/cron.log 2>&1
